# deploy/inference-console.service
[Unit]
Description=Inference Console (serveur web)
After=network.target

[Service]
Type=simple
# Adapte WorkingDirectory au chemin de déploiement sur la VM.
WorkingDirectory=/opt/hackathon_ia_console
ExecStart=/usr/bin/python3 /opt/hackathon_ia_console/server.py
# Variables : adapte le PORT à celui attendu par la passerelle klouders
# (https://hack-ia01.klouders.fr/ -> cette VM).
Environment=HOST=0.0.0.0
Environment=PORT=8080
# Proxy interne vers le backend d'inférence (évite le CORS).
Environment=UPSTREAM=http://localhost:11434
Restart=on-failure
RestartSec=3
# Durcissement léger.
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
